# solaris_oci/files/identity.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
 Copyright (c) 2004, 2018, Oracle and/or its affiliates. All rights reserved.

	NOTE:  This service manifest is not editable; its contents will
	be overwritten by package or patch operations, including
	operating system upgrade.  Make customizations in a different
	file.
-->

<service_bundle type='manifest' name='SUNWcsr:identity'>

<service
	name='system/identity'
	type='service'
	version='1'>

	<exec_method
		type='method'
		name='stop'
		exec=':true'
		timeout_seconds='60' />

	<property_group name='startd' type='framework'>
		<propval name='duration' type='astring' value='transient' />
	</property_group>

	<instance name='node' enabled='true'>
		<dependent
			name='identitynode_self-assembly-complete'
			grouping='optional_all'
			restart_on='none'>
			<service_fmri value='svc:/milestone/self-assembly-complete' />
		</dependent>

		<exec_method
			type='method'
			name='start'
			exec='/lib/svc/method/identity-node start'
			timeout_seconds='60' />

		<exec_method
			type='method'
			name='refresh'
			exec='/lib/svc/method/identity-node refresh'
			timeout_seconds='60' />

		<exec_method
			type='method'
			name='unconfigure'
			exec='/lib/svc/method/identity-node unconfigure'
			timeout_seconds='60' />

		<property_group name='config' type='application'>
			<propval name='nodename' type='astring' value=''/>
			<propval name='loopback' type='astring' value=''/>
			<propval name='enable_mapping' type='boolean'
			    value='true'/>
			<propval name='ignore_dhcp_hostname' type='boolean'
			    value='false'/>
		</property_group>

		<property_group name='private_data' type='application'>
			<propval name='mapped_hostname' type='astring'
			    value=''/>
			<property name='mapped_addresses' type='net_address'>
				<net_address_list>
					<value_node value='0.0.0.0' />
				</net_address_list>
			</property>
		</property_group>

		<property_group name='sysconfig' type='sysconfig'>
			<propval name='group' type='astring' value='identity' />
			<propval name='reconfigurable' type='boolean' value='false'/>
			<property name='config_properties' type='astring'>
				<astring_list>
					<value_node value='sc_hostname:config/nodename' />
				</astring_list>
			</property>
		</property_group>

		<template>
			<common_name>
				<loctext xml:lang='C'>
					system identity (nodename)
				</loctext>
			</common_name>
			<documentation>
				<manpage title='nodename' section='5'
					manpath='/usr/share/man' />
			</documentation>
			<pg_pattern name="config" type="application">
				<description>
					<loctext xml:lang='C'>
Properties that control the configuration of the node
					</loctext>
				</description>
				<prop_pattern name="nodename" type="astring">
					<description>
						<loctext xml:lang='C'>
Network name of the computer
						</loctext>
					</description>
				</prop_pattern>
				<prop_pattern name="loopback" type="astring">
					<description>
						<loctext xml:lang='C'>
Hostname mapped to loopback
						</loctext>
					</description>
				</prop_pattern>
				<prop_pattern name="ignore_dhcp_hostname"
				    type="boolean">
					<description>
						<loctext xml:lang='C'>
Ignore hostname returned by DHCP
						</loctext>
					</description>
				</prop_pattern>
			</pg_pattern>
		</template>

	</instance>

	<stability value='Unstable' />

</service>

</service_bundle>
